# stockflow-telegram-backup.service
[Unit]
Description=Stock Flow daily Telegram backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/var/www/stock_flow
EnvironmentFile=/var/www/stock_flow/.env
Environment=TELEGRAM_BOT_TOKEN=8000825371:AAHpWNnggP_DrtM-zp21_O6Xt0g80wdWtwc
Environment=TELEGRAM_STORAGE_CHAT_ID=-1003752537433
ExecStart=/var/www/stock_flow/venv/bin/python /var/www/stock_flow/stockflow_telegram_backup.py
